# randomWallpaper.cron
*/15 * * * * DISPLAY=:0.0 ~/work/github/AlfScripts/randomWallpaper.sh
